# infrastructure/ansible/roles/health-checks/templates/ablage-health-check.timer.j2
[Unit]
Description=Ablage-System Health Check Timer
Documentation=https://github.com/org/ablage-system

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=30
Persistent=true

[Install]
WantedBy=timers.target
